# mqtt-io/rootfs/etc/s6-overlay/s6-rc.d/mqtt-io/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: MQTT IO
# Runs MQTT IO
# ==============================================================================
declare configuration_file

configuration_file=$(bashio::config 'configuration_file')
if ! bashio::fs.file_exists "${configuration_file}"; then
    bashio::log.fatal
    bashio::log.fatal "Seems like the configured configuration file does"
    bashio::log.fatal "not exists:"
    bashio::log.fatal
    bashio::log.fatal "${configuration_file}"
    bashio::log.fatal
    bashio::log.fatal "Please check the add-on configuration, or create the"
    bashio::log.fatal "configuration file. More information about the"
    bashio::log.fatal "configuration file format see:"
    bashio::log.fatal
    bashio::log.fatal "https://mqtt-io.app"
    bashio::log.fatal
    bashio::exit.nok
fi

bashio::log.info "Starting MQTT IO..."
exec python3 -m mqtt_io "$(bashio::config 'configuration_file')"
